# scripts/hypertls-microservice.service
[Service]
ExecStart=/home/lzuccarelli/Projects/microservices/hypertls-microservice /home/lzuccarelli/Projects/microservices
